# files/upstart/app.conf
description "App"

start on runlevel [2345]
stop on runlevel [016]

pre-start script
  mkdir -p /home/escritoriovirtual/log
  chown -R escritoriovirtual. /home/escritoriovirtual/log
end script
